# crontab
1 1 * * * python /app/flick_email.py -v
